# home/.config/systemd/user/wallpaper-changer.service
[Unit]
Description=Change wallpaper

[Service]
Type=oneshot
ExecStart=/home/%u/.config/i3/bin/wallpaper-changer

[Install]
WantedBy=default.target
